# containers/jmx-exporter/1/data/profile.d/jmx_exporter.sh
export PATH="/opt/drycc/jmx_exporter/bin:$PATH"
